# supervisord.conf
[supervisord]
nodaemon=true

[program:socket-server]
command=/app/socket-server

[program:nginx]
command=nginx -g "daemon off;"
